# board/crystalfontz/cfa92x/root/etc/rc.local
#modprobe asix
modprobe can
modprobe can-raw
modprobe flexcan
ip link set can0 type can bitrate 500000
ifconfig can0 up
ntpdate 192.168.1.1
